# scripts/trading-morning.service
[Unit]
Description=Morning Trading Automation Service
After=network.target

[Service]
Type=oneshot
# Don't specify User= in user services - it runs as the user who started it
WorkingDirectory=/home/taiconley/Desktop/Projects/trading
# Use the virtual environment's Python interpreter
ExecStart=/home/taiconley/Desktop/Projects/trading/env/bin/python3 /home/taiconley/Desktop/Projects/trading/scripts/morning_automation.py
StandardOutput=append:/home/taiconley/Desktop/Projects/trading/logs/morning_automation.log
StandardError=append:/home/taiconley/Desktop/Projects/trading/logs/morning_automation.log
# Don't kill child processes (like TWS) when the service exits
KillMode=process

# Environment variables for X11 access
Environment="DISPLAY=:1"
Environment="XAUTHORITY=/run/user/1000/gdm/Xauthority"
Environment="XDG_RUNTIME_DIR=/run/user/1000"
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus"
Environment="PATH=/home/taiconley/Desktop/Projects/trading/env/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="VIRTUAL_ENV=/home/taiconley/Desktop/Projects/trading/env"

# Restart policy
Restart=no

[Install]
WantedBy=default.target

